# frp_client_mtls/rootfs/etc/services.d/frpc/run
#!/usr/bin/env bashio
# shellcheck shell=bash

CONFIG=/data/frpc.ini

/usr/local/bin/frpc -c /data/frpc.ini
